# glitch-core.service
[Unit]
Description=GLITCH.CORE SSH Daemon
After=network.target

[Service]
ExecStart=/usr/bin/npx tsx /opt/glitch-core/daemon.ts
WorkingDirectory=/opt/glitch-core
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
Restart=always
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
